# benchmark_demo/webserver_app/linux_app/rpmsg_json.service
[Unit]
Wants=network-online.target
After=network.target network-online.target

[Service]
EnvironmentFile=/etc/benchmark_server.conf
ExecStart=rpmsg_json ${APP_DIR}/oob_data.json

[Install]
WantedBy=multi-user.target
