# benchmark_demo/webserver_app/webserver/benchmark_server.service
[Unit]
Wants=network-online.target
After=network.target network-online.target

[Service]
EnvironmentFile=/etc/benchmark_server.conf
ExecStart=benchmark_server $APP_DIR

[Install]
WantedBy=multi-user.target
